# crontab
0,30 8-20 * * mon-fri ~/.vim/journal.sh
